# backend/systemd/fpl-livefpl-predictions.timer
[Unit]
Description=Run LiveFPL predictions scrape every 30 minutes

[Timer]
OnCalendar=*:00/30
OnBootSec=2min
Persistent=yes

[Install]
WantedBy=timers.target
